# server_files/supervisor.conf
; Supervisor program for the storefront. Replace placeholder paths/user.

[program:myproject]
command = /srv/myproject/server_files/gunicorn_start
user = www-data
stdout_logfile = /srv/myproject/logs/supervisor.log
redirect_stderr = true
environment=LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8
